# debian/postinst
#!/bin/bash

set -eo pipefail

case "${1}" in
    configure|upgrade|reinstall     )
        systemctl is-active --quiet ddnsd.service || systemctl enable --now ddnsd.service
    ;;
    *                               )
        echo "postinst called with unknown argument: $1" >&2 ; exit 1
    ;;
esac
